# config_examples/email-worker.service
[Unit]
Description=Email worker

[Service]
WorkingDirectory=/path/to/app/task-tracker
ExecStart=/usr/bin/npm run email-worker