# config_examples/email-worker.service
[Unit]
Description=Email worker

[Service]
WorkingDirectory=/home/ubuntu/_projects/task-tracker-blue
ExecStart=/usr/bin/npm run email-worker